# data/autostart/pcmanfm.desktop
[Desktop Entry]
Version=1.0
Name=PCManFM
Comment=PCManFM
Exec=pcmanfm -d
X-Vera-Autostart-Phase=Desktop
