# data/autostart/volumeicon.desktop
[Desktop Entry]
Version=1.0
Name=volumeicon
Comment=Volumeicon
Exec=volumeicon
X-Vera-Autostart-Phase=Other
